# rootfs/etc/s6-overlay/s6-rc.d/tar1090/run
#!/bin/sh
exec /etc/s6-overlay/scripts/tar1090
