# rootfs/etc/s6-overlay/s6-rc.d/readsb/finish
#!/bin/sh
rm -rf /run/readsb
rm -f /run/readsb-prometheus.prom
