# rootfs/etc/s6-overlay/s6-rc.d/graphs1090-writeback/run
#!/bin/sh
exec /etc/s6-overlay/scripts/graphs1090-writeback
