# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

mkdir -p /var/log/nginx

exec \
  /usr/sbin/nginx 2>&1 | mawk -W interactive '{print "[" strftime("%Y-%m-%d %H:%M:%S", systime()) "][nginx] " $0}'
