# rootfs/etc/services.d/collectd/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash disable=SC1091

source /scripts/common

if chk_enabled "${GRAPHS1090_DISABLE}"; then
    exec sleep infinity
fi

# wait for neccessary input file to exist before starting

# shellcheck disable=SC2016
s6wrap --prepend=collectd --timestamps --args echo "Waiting for readsb to start..."
while [[ ! -e /usr/share/graphs1090/data-symlink/data/stats.json ]]; do
    sleep 1
done

# test config?
/usr/sbin/collectd -t -T -f -C /etc/collectd/collectd.conf > /dev/null 2>&1

#shellcheck disable=SC2016
exec s6wrap --prepend=collectd --timestamps --args /usr/sbin/collectd -C /etc/collectd/collectd.conf -f


