# rootfs/etc/s6-overlay/s6-rc.d/03-temp-fix/up
#!/bin/sh
exec /etc/s6-overlay/s6-rc.d/scripts/03-temp-fix-up
