# rootfs/etc/s6-overlay/s6-rc.d/04-tar1090-configure/up
#!/bin/sh
exec /etc/s6-overlay/scripts/04-tar1090-configure
