# rootfs/etc/s6-overlay/s6-rc.d/08-graphs1090/up
#!/bin/sh
exec /etc/s6-overlay/scripts/08-graphs1090-init
