# rootfs/etc/s6-overlay/s6-rc.d/timelapse1090/run
#!/bin/sh
exec /etc/s6-overlay/scripts/timelapse1090
